# etc/supervisord.conf
[unix_http_server]
file=%(here)s/../var/supervisor.sock

[supervisord]
logfile=%(here)s/../var/supervisord.log
logfile_maxbytes=10MB
logfile_backups=20 
loglevel=info
pidfile=%(here)s/../var/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/../var/supervisor.sock

[program:voteit]
command = %(here)s/../bin/paster serve %(here)s/../etc/production.ini
redirect_stderr = true
stdout_logfile=%(here)s/../var/log/%(program_name)s-%(process_num)01d.log
stdout_capture_maxbytes=10MB
stdout_logfile_backups=20

[program:zeo]
command = %(here)s/../bin/runzeo -C %(here)s/../etc/zeo.conf
redirect_stderr = true
stdout_capture_maxbytes=10MB
stdout_logfile_backups=20
